# l/ltt-control/lttng-tools_2.9.3-1_amd64/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        if which addgroup >/dev/null; then
            addgroup --system --quiet tracing > /dev/null || true
        fi
    ;;
esac

if dpkg-maintscript-helper supports rm_conffile ; then
    dpkg-maintscript-helper rm_conffile \
        /etc/bash_completion.d/lttng 2.8.2-1~ lttng-tools -- "$@"
fi

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

# Automatically added by dh_installinit
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ]; then
	if [ -x "/etc/init.d/lttng-sessiond" ]; then
		update-rc.d lttng-sessiond defaults >/dev/null
		invoke-rc.d lttng-sessiond start || exit $?
	fi
fi
# End automatically added section


exit 0
